# services/backups/crontab
# Run backup script every day at midnight
0 0 * * * /app/cron.sh >>/var/log/cron.log 2>&1
